# etc/init.d/elogind-openrc
#!/sbin/openrc-run

description="Starts the elogind user, seat, and session manager"
command="/usr/libexec/elogind"
command_args="--no-daemon"
supervisor=supervise-daemon

depend() {
	need dbus
	provide elogind
}
